# supervisord.conf
[supervisord]
nodaemon=true

[program:fastapi]
command=uvicorn main:app --host 0.0.0.0 --port 8000

[program:gradio]
command=python3 gradio_app.py